# prep/prep_nextcloud_start.service
[Unit]
After=network.service

[Service]
ExecStart=/root/prep-nextcloud-start.sh

[Install]
WantedBy=default.target
